# deploy/examples/makereadyos-automations.timer
[Unit]
Description=Run MakeReadyOS scheduled automation evaluation hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true
Unit=makereadyos-automations.service

[Install]
WantedBy=timers.target
