# deploy/examples/makereadyos-backup.service
[Unit]
Description=Create and prune a full MakeReadyOS backup snapshot
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=makereadyos
Group=docker
WorkingDirectory=/opt/makereadyos
ExecStart=/opt/makereadyos/backup-all.sh

# The script reads /opt/makereadyos/.env. Keep credentials there, not in this unit.
